# systemd/kite-eod-session-plan.timer
[Unit]
Description=EOD session plan at 19:00 IST on weekdays (after NSE FII/DII flow publish ~18:30 IST)

[Timer]
OnCalendar=Mon-Fri *-*-* 19:00:00 Asia/Kolkata
Persistent=true

[Install]
WantedBy=timers.target
