# systemd/kite-fao-reports.service
[Unit]
Description=NSE F&O daily reports pull + pre-market bias (full bundle, one-shot)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/root/eodreport
Environment=FAO_ALL=1
ExecStart=/usr/bin/env bash /root/eodreport/cron_fao_reports.sh
StandardOutput=append:/root/eodreport/fao-reports.log
StandardError=append:/root/eodreport/fao-reports.log
