# systemd/kite-fao-reports.service
[Unit]
Description=NSE F&O daily reports pull + pre-market bias (full bundle, one-shot)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/root/eodreport
EnvironmentFile=-/root/eodreport/.env
Environment=FAO_ALL=1
ExecStart=/usr/bin/env bash /root/eodreport/cron_fao_reports.sh
ExecStartPost=-/root/eodreport/pull_latest_reports.sh --fao --sync
StandardOutput=append:/root/eodreport/fao-reports.log
StandardError=append:/root/eodreport/fao-reports.log
